# src/scripts/busybee.service
[Unit]
Description=BusyBee - Universal Calendar Sync Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jimi
WorkingDirectory=/home/jimi/devP/home/busybee
ExecStart=/home/jimi/devP/home/busybee/scripts/busybee.sh
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
